# stow/launchd/Library/LaunchAgents/com.rzman.pihole-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rzman.pihole-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rzman/.bin/pihole-backup</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/rzman/.aria2/pihole-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rzman/.aria2/pihole-backup.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/rzman</string>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
